# supervisor/waze_scheduler.conf
; =============================================================================
; Supervisor — Produção (Hostinger)
; Copie para: /etc/supervisor/conf.d/wazebr.conf
; Após alterações: supervisorctl reread && supervisorctl update
; =============================================================================

; -----------------------------------------------------------------------------
; Scheduler — dispara as mensagens agendadas para as filas
; 1 processo é suficiente (só agenda, não processa)
; -----------------------------------------------------------------------------
[program:waze_scheduler]
command=php /home/u629736858/domains/acheireviews.com.br/public_html/trafik/bin/console messenger:consume scheduler_waze_feed --time-limit=3600 --memory-limit=64M
directory=/home/u629736858/domains/acheireviews.com.br/public_html/trafik
autostart=true
autorestart=true
startretries=5
stopwaitsecs=30
stdout_logfile=/home/u629736858/logs/supervisor/waze_scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/home/u629736858/logs/supervisor/waze_scheduler_error.log
redirect_stderr=false
environment=APP_ENV="prod",APP_DEBUG="0"

; -----------------------------------------------------------------------------
; Workers Waze — processam alertas, jams, feeds PartnerHub/TVT em paralelo
; numprocs=2 → 2 processos simultâneos → FetchAlerts e FetchTraffic ao mesmo tempo
; -----------------------------------------------------------------------------
[program:worker_waze]
command=php /home/u629736858/domains/acheireviews.com.br/public_html/trafik/bin/console messenger:consume async_waze --time-limit=3600 --memory-limit=128M
directory=/home/u629736858/domains/acheireviews.com.br/public_html/trafik
autostart=true
autorestart=true
startretries=5
stopwaitsecs=60
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/home/u629736858/logs/supervisor/worker_waze_%(process_num)02d.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/home/u629736858/logs/supervisor/worker_waze_%(process_num)02d_error.log
redirect_stderr=false
environment=APP_ENV="prod",APP_DEBUG="0"

; -----------------------------------------------------------------------------
; Worker CEMADEN — processa dados hidrológicos (menor frequência)
; numprocs=1 — 1 processo é suficiente (coleta a cada 10 min)
; -----------------------------------------------------------------------------
[program:worker_cemaden]
command=php /home/u629736858/domains/acheireviews.com.br/public_html/trafik/bin/console messenger:consume async_cemaden --time-limit=3600 --memory-limit=64M
directory=/home/u629736858/domains/acheireviews.com.br/public_html/trafik
autostart=true
autorestart=true
startretries=5
stopwaitsecs=60
numprocs=1
stdout_logfile=/home/u629736858/logs/supervisor/worker_cemaden.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/home/u629736858/logs/supervisor/worker_cemaden_error.log
redirect_stderr=false
environment=APP_ENV="prod",APP_DEBUG="0"
