# supervisor/waze_scheduler_dev.conf
; =============================================================================
; Supervisor — Desenvolvimento local
; =============================================================================

[program:waze_scheduler_dev]
command=php /var/www/wazeBR-symfony/bin/console messenger:consume scheduler_waze_feed --time-limit=3600 --memory-limit=64M -vv
directory=/var/www/wazeBR-symfony
autostart=true
autorestart=true
startretries=3
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/waze_scheduler_dev.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/waze_scheduler_dev_error.log
redirect_stderr=false
environment=APP_ENV="dev",APP_DEBUG="1"

[program:worker_waze_dev]
command=php /var/www/wazeBR-symfony/bin/console messenger:consume async_waze --time-limit=3600 --memory-limit=128M -vv
directory=/var/www/wazeBR-symfony
autostart=true
autorestart=true
startretries=3
stopwaitsecs=60
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor/worker_waze_dev_%(process_num)02d.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/worker_waze_dev_%(process_num)02d_error.log
redirect_stderr=false
environment=APP_ENV="dev",APP_DEBUG="1"

[program:worker_cemaden_dev]
command=php /var/www/wazeBR-symfony/bin/console messenger:consume async_cemaden --time-limit=3600 --memory-limit=64M -vv
directory=/var/www/wazeBR-symfony
autostart=true
autorestart=true
startretries=3
stopwaitsecs=60
stdout_logfile=/var/log/supervisor/worker_cemaden_dev.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/worker_cemaden_dev_error.log
redirect_stderr=false
environment=APP_ENV="dev",APP_DEBUG="1"
